# backup.conf
[program:backup]
command=python /home/pi/PiCameraRig/camera-node/backup.py
user=pi
autostart=true
autorestart=true
