# camera.conf
[program:camera]
command=python /home/pi/PiCameraRig/camera-node/run.py
user=pi
autostart=true
autorestart=true
