# controller/server.conf
[program:server]
directory=/home/pi/controller/
command=python -m SimpleHTTPServer
user=pi
autostart=true
autorestart=true
